# .config/systemd/user/obsidian-backup.timer
[Unit]
Description=Run obsidian-backup daily and on boot

[Service]
OnBootSec=1m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
