# .config/systemd/user/obsidian-backup.timer
[Unit]
Description=Run obsidian-backup hourly and on boot

[Timer]
OnBootSec=1m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
